# deploy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec nginx -c /config/nginx/nginx.conf
